# miniweb/root/etc/init.d/autoupgrade
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=103

start() {
		/sbin/upgrade_rule &>/dev/null	
}
